# apache/templates/auth_logger.service.j2
[Unit]
Description="Log all apache auth credential attempts and successful logins"
After=network.target

[Service]
ExecStart=/usr/local/bin/auth_logger -i {{ auth_logger_sniffer_iface }}
